# deploy/systemd/constellate-agent.service
[Unit]
Description=Constellate agent (connect relay)
After=network-online.target constellate-session-host.service
Wants=network-online.target
Requires=constellate-session-host.service

[Service]
Type=simple
# Replace YOUR_USER with the service account, or let `constellate-agent install`
# generate a populated unit for you.
User=YOUR_USER
ExecStart=/usr/local/bin/constellate-agent connect
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
